# ids-dock-tracker.service
[Unit]
Description=IDS Dock Tracker Service
After=network.target postgresql.service

[Service]
Type=simple
User=docktracker
Environment=ENVIRONMENT=production
Environment=LOG_LEVEL=info
ExecStart=/usr/local/bin/dock-tracker
Restart=always
RestartSec=5

# Enable logging
StandardOutput=append:/var/log/dock-tracker/application.log
StandardError=append:/var/log/dock-tracker/error.log

[Install]
WantedBy=multi-user.target
